# deploy/systemd/cloudflared.service
[Unit]
Description=Cloudflare Tunnel for Madani Muzakarah
After=network.target nginx.service

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --config /etc/cloudflared/config.yml run
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
